# AcerMonitor_Task.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Importar no Task Scheduler do Windows para iniciar o watchdog automaticamente.
  O watchdog por sua vez gerencia o monitor.py, reiniciando-o se necessário.

  Como importar:
  1. Abra o Task Scheduler (Agendador de Tarefas)
  2. Clique em "Import Task..." no painel direito
  3. Selecione este arquivo .xml
  4. Ajuste o caminho em <WorkingDirectory> para o local correto
  5. Clique OK
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Acer Aspire 5 — Watchdog do Monitor de Reinicializações</Description>
    <Author>AcerMonitor</Author>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>999</Count>
    </RestartOnFailure>
  </Settings>
  <Actions>
    <Exec>
      <!--
        AJUSTE ESTE CAMINHO:
        Troque C:\Users\SEU_USUARIO\acer-crash-monitor pelo caminho real da pasta
      -->
      <Command>pythonw.exe</Command>
      <Arguments>watchdog.py</Arguments>
      <WorkingDirectory>C:\Users\SEU_USUARIO\acer-crash-monitor</WorkingDirectory>
    </Exec>
  </Actions>
</Task>